# systemd/proxmox-monitor.service
[Unit]
Description=Proxmox VE Telegram Monitor
Documentation=https://github.com/yourusername/proxmox-monitor
After=network.target
Wants=proxmox-monitor.timer

[Service]
Type=simple
User=root
WorkingDirectory=/opt/proxmox-monitor

# Set PATH to include root's local bin for uv
Environment="PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

ExecStart=/bin/sh /opt/proxmox-monitor/scripts/entrypoint.sh

# Restart policy
Restart=always
RestartSec=30

# Resource limits
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxmox-monitor

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/proxmox-monitor /var/log

[Install]
WantedBy=multi-user.target
